# llm-enrichment.service

[Unit]
Description=LLM Alert Enrichment Service
After=network.target


[Service]
Type=simple
User=youruser
WorkingDirectory=/path/to/your/project
EnvironmentFile=/path/to/your/project/.env
ExecStart=/usr/bin/python3 enrichment_service.py start
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal


[Install]
WantedBy=multi-user.target
